# example/credits.node.systemd.service.in
[Unit]
Description=credits.node%i service
After=network-online.target syslog.target

[Service]
LimitCORE=infinity
Type=forking
EnvironmentFile=${ENV_FILE}
User=${USER}
Group=${GROUP}
ExecStart=${CREDITSHOME}/bin/node
PIDFile=${CREDITSHOME}/node.%i.pid
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
